# debian/postinst
#!/bin/bash
# postinst script for internet-monitor

set -e

# Source debconf library
. /usr/share/debconf/confmodule

case "$1" in
    configure)
        # Set ownership on installed sound files
        SOUND_DIR="/usr/share/asterisk/sounds/custom"
        if getent passwd asterisk >/dev/null 2>&1; then
            [ -f "${SOUND_DIR}/internet-yes.ulaw" ] && chown asterisk:asterisk "${SOUND_DIR}/internet-yes.ulaw" 2>/dev/null || true
            [ -f "${SOUND_DIR}/internet-no.ulaw" ] && chown asterisk:asterisk "${SOUND_DIR}/internet-no.ulaw" 2>/dev/null || true
        fi
        
        # Reload systemd daemon
        if command -v systemctl >/dev/null 2>&1; then
            systemctl daemon-reload || true
        fi
        
        # Install default configuration if it doesn't exist
        if [ -f /etc/internet-monitor.conf.dpkg-new ]; then
            if [ ! -f /etc/internet-monitor.conf ]; then
                mv /etc/internet-monitor.conf.dpkg-new /etc/internet-monitor.conf
                echo "Created default configuration file at /etc/internet-monitor.conf"
            else
                rm -f /etc/internet-monitor.conf.dpkg-new
            fi
        fi
        
        # Check if configuration needs to be set up
        if [ ! -f /etc/internet-monitor.conf ] || ! grep -q "^NODE_NUMBER=" /etc/internet-monitor.conf 2>/dev/null || grep -q "^NODE_NUMBER=12345$" /etc/internet-monitor.conf 2>/dev/null; then
            echo ""
            echo "=========================================="
            echo "Internet Monitor Configuration Required"
            echo "=========================================="
            echo ""
            echo "Please edit /etc/internet-monitor.conf and set your NODE_NUMBER"
            echo "Then start the service with:"
            echo "  sudo systemctl start internet-monitor"
            echo "  sudo systemctl enable internet-monitor"
            echo ""
        fi
        ;;
    
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    
    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
